# data_record.conf
[program:data_record]
directory=/home/c8763yee/Grafana_DashBoard_v4
autostart=true
autorestart=true
startretries=10
user=c8763yee
logfile_maxbytes=0
logfile_backups=0
command=python3 app.py
